# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

/usr/sbin/cron -f -S -l 0 -c /etc/crontab
